# files/system/etc/profile.d/99-a11y.sh
export GTK_MODULES="atk-bridge"
export QT_ACCESSIBILITY=1
export ENABLE_ACCESSIBILITY=1