# files/system/etc/xdg/autostart/orca.desktop
[Desktop Entry]
Type=Application
Name=Orca
Exec=orca
X-GNOME-Autostart-enabled=true
X-KDE-autostart-after=panel